# wazo-asterisk-restart.timer
[Unit]
Description=Daily restart of Asterisk

[Timer]
OnCalendar=*-*-* 23:30:00
Persistent=true

[Install]
WantedBy=timers.target
